# AltServer.service.dist
[Unit]
Description=AltServer

[Service]
ExecStart=<path>/server.sh

[Install]
WantedBy=multi-user.target