# services/AltServer.service.dist
[Unit]
Description=AltServer

[Service]
ExecStart=<path>/AltServer

[Install]
WantedBy=multi-user.target