# services/netmuxd.service.dist
[Unit]
Description=netmuxd

[Service]
ExecStart=<path>/netmuxd

[Install]
WantedBy=multi-user.target